# supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
#command=/usr/sbin/php-fpm8.3 #ubuntu
command=/usr/local/sbin/php-fpm #php

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
